# terminal-server.service
[Unit]
Description=WebSocket Terminal Server for Terminal Launcher
After=network.target

[Service]
Type=simple
User=cactus
WorkingDirectory=/home/cactus/claude/terminal-launcher
ExecStart=/usr/bin/python3 /home/cactus/claude/terminal-launcher/terminal-server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Resource limits
LimitNOFILE=65536

# Security (allow PTY operations)
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
